# .config/systemd/user/recall.timer
[Unit]
Description=Timer für RECALL

[Timer]
OnCalendar=*-*-* *:*:0/10

[Install]
WantedBy=timers.target
