# usr/systemd/user/recall.service
[Unit]
DefaultDependencies=no
Description=Linux am Dienstag RECALL Experiment

[Service]
Type=oneshot
ExecStart=/usr/sbin/recall
